# crates/famp/tests/fixtures/sample-com.famp.broker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.famp.broker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/.cargo/bin/famp</string>
        <string>broker</string>
        <string>--no-idle-exit</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/.famp/broker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/.famp/broker.log</string>
</dict>
</plist>
